# files/system/etc/systemd/system/nvidia-settings-toggle.service
[Unit]
Description=Toggle NVIDIA Settings visibility based on GPU presence
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/toggle-nvidia-settings.sh
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
